# settings/inori.ini
[program:inori]
command=/srv/virtualenvs/inorienv/bin/gunicorn -w 4 inori.app:app -b 0.0.0.0:17020 --preload
autostart=true
autorestart=unexpected
startsecs=1
startretries=3
exitcodes=0,2
stopasgroup=false
killasgroup=false
redirect_stderr=true
stdout_logfile=/var/log/inori/gunicorn.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=10
stderr_logfile=/var/log/inori/gunicorn.err
stderr_logfile_maxbytes=1MB
stderr_logfile_backups=10
environment=PYTHONPATH=/srv/inori
